# roles/osquery/files/osquery.conf
{
    "options": {
        "disable_events": false,
        "enable_file_events": true,
        "enable_bpf_events": true,
        "schedule_splay_percent": 10,
        "events_expiry": 3600,
        "schedule_default_interval": 3600
    },
    "platform": "linux",
    "schedule": {
        "suspicious_commands": {
            "query": "select * from bpf_process_events where path in ('/bin/bash', '/bin/csh', '/bin/dash', '/bin/busybox', '/bin/zsh', '/usr/bin/which', '/usr/bin/curl', '/usr/bin/wget', '/usr/bin/tar', '/usr/bin/base64', '/usr/bin/chmod', '/usr/bin/uname', '/usr/bin/whoami', '/usr/bin/id', '/usr/sbin/ifconfig', '/usr/bin/netstat','/usr/bin/lsof', '/usr/bin/telnet', '/usr/bin/nc', '/usr/sbin/route', '/usr/bin/ping', '/usr/bin/ssh', '/usr/bin/ftp');",
            "interval": 60,
            "description": "Detect suspicious commands like shell, enumeration, reconnaissance, lateral movement, persistance"
        },
        "blacklist_process_events": {
            "query": "select * from bpf_process_events where uid in (33);",
            "interval": 60,
            "description": "Looking for command execution for www-data(33)"
        },
        "blacklist_socket_events": {
            "query": "select * from bpf_socket_events where uid in (33) and syscall='connect' and remote_address not in ('127.0.0.1', '127.0.0.53') and remote_port>0;",
            "interval": 60,
            "description": "Looking for initiated connections by www-data(33)"
        },
        "file_events": {
            "query": "SELECT * FROM file_events;",
            "interval": 60,
            "description": "File events collected from file integrity monitoring",
            "removed": false
        },
        "ld_preload": {
            "query": "SELECT process_envs.pid, process_envs.key, process_envs.value, processes.name, processes.path, processes.cmdline, processes.cwd FROM process_envs join processes USING (pid) WHERE key = 'LD_PRELOAD';",
            "interval": 60,
            "description": "Any processes that run with an LD_PRELOAD environment variable"
        },
        "usb_devices": {
            "query": "SELECT * FROM usb_devices;",
            "interval": 120,
            "description": "Retrieves the current list of USB devices in the target system."
        },
        "behavioral_reverse_shell": {
            "query": "SELECT DISTINCT(processes.pid), processes.parent, processes.name, processes.path, processes.cmdline, processes.cwd, processes.root, processes.uid, processes.gid, processes.start_time, process_open_sockets.remote_address, process_open_sockets.remote_port, (SELECT cmdline FROM processes AS parent_cmdline WHERE pid=processes.parent) AS parent_cmdline FROM processes JOIN process_open_sockets USING (pid) LEFT OUTER JOIN process_open_files ON processes.pid = process_open_files.pid WHERE (name='sh' OR name='bash') AND remote_address NOT IN ('0.0.0.0', '::', '') AND remote_address NOT LIKE '127.0.0.%';",
            "interval": 600,
            "description": "Find shell processes that have open sockets"
        },
        "iptables": {
            "query": "SELECT * FROM iptables;",
            "interval": 600,
            "description": "Retrieves the current filters and chains per filter in the target system."
        },
        "network_interfaces": {
            "query": "SELECT a.interface, a.address, d.mac FROM interface_addresses a JOIN interface_details d USING (interface);",
            "interval": 600,
            "description": "Record the network interfaces and their associated IP and MAC addresses"
        },
        "authorized_keys": {
            "query": "SELECT * FROM users CROSS JOIN authorized_keys USING (uid);",
            "interval": 3600,
            "description": "A line-delimited authorized_keys table."
        },
        "crontab": {
            "query": "SELECT * FROM crontab;",
            "interval": 3600,
            "description": "Retrieves all the jobs scheduled in crontab in the target system."
        },
        "startup_items": {
            "query": "select * from startup_items where status='active';",
            "interval": 3600,
            "description": "Retrieves all active startup items in the target system."
        },
        "dns_resolvers": {
            "query": "SELECT * FROM dns_resolvers;",
            "interval": 3600,
            "description": "DNS resolvers used by the host"
        },
        "etc_hosts": {
            "query": "SELECT * FROM etc_hosts;",
            "interval": 3600,
            "description": "Retrieves all the entries in the target system /etc/hosts file."
        },
        "kernel_modules": {
            "query": "SELECT * FROM kernel_modules;",
            "interval": 3600,
            "description": "Linux kernel modules both loaded and within the load search path."
        },
        "last": {
            "query": "SELECT * FROM last;",
            "interval": 3600,
            "description": "Retrieves the list of the latest logins with PID, username and timestamp."
        },
        "ld_so_preload_exists": {
            "query": "SELECT * FROM file WHERE path='/etc/ld.so.preload' AND path!='';",
            "interval": 3600,
            "description": "Generates an event if ld.so.preload is present - used by rootkits such as Jynx"
        },
        "listening_ports": {
            "query": "SELECT pid, port, processes.path, cmdline, cwd FROM listening_ports JOIN processes USING (pid) WHERE port!=0;",
            "interval": 3600,
            "description": "Gather information about processes that are listening on a socket.",
            "removed": false
        },
        "system_running_processes": {
            "query": "select pr.pid, pr.name, usr.username, pr.path, pr.cmdline from processes pr LEFT JOIN users usr ON pr.uid = usr.uid WHERE pr.cmdline != '';",
            "interval": 3600,
            "description": "List Linux System running processes with CMDLINE not null. - ATT&CK T1059,T1108,T1166,T1100,T1064,T1107,T1003,T1033,T1016,T1082,T1057,T1201,T1083,T1217,T1087,T1072,T1002"
        },
        "process_open_sockets": {
            "query": "select DISTINCT p.name, p.path, pos.remote_address, pos.remote_port from process_open_sockets pos LEFT JOIN processes p ON pos.pid = p.pid WHERE pos.remote_port != 0 AND p.name != '';",
            "interval": 3600,
            "description": "Returns the network connections from system processes - ATT&CK T1108,T1100,T1102,T1105,T1039,T1029,T1011,T1041,T1043,T1090,T1094,T1048,T1132,T1020,T1065,T1001,T1071,T1219,T1104,T1008",
            "removed": false
        },
        "shell_history": {
            "query": "SELECT * FROM users CROSS JOIN shell_history USING (uid);",
            "interval": 3600,
            "description": "Record shell history for all users on system (instead of just root)"
        },
        "sudoers": {
            "query": "select * from sudoers;",
            "interval": 3600,
            "description": "Linux sudoers information - ATT&CK T1169,T1206"
        },
        "suid_bin": {
            "query": "SELECT * FROM suid_bin;",
            "interval": 86400,
            "description": "Display any SUID binaries that are owned by root"
        },
        "user_ssh_keys": {
            "query": "SELECT * FROM users CROSS JOIN user_ssh_keys USING (uid);",
            "interval": 86400,
            "description": "Returns the private keys in the users ~/.ssh directory and whether or not they are encrypted"
        },
        "users": {
            "query": "SELECT * FROM users;",
            "interval": 86400,
            "description": "Local system users."
        },
        "user_groups": {
            "query": "select ug.uid, ug.gid, u.username, g.groupname from user_groups ug left join users u on ug.uid=u.gid left join groups g on ug.gid=g.gid;",
            "interval": 86400,
            "description": "Local system user_groups."
        },
        "mounts": {
            "query": "SELECT device, device_alias, path, type, blocks_size, flags FROM mounts;",
            "interval": 86400,
            "description": "Retrieves the current list of mounted drives in the target system."
        },
        "deb_packages": {
            "query": "SELECT * FROM deb_packages;",
            "interval": 86400,
            "description": "Display all installed DEB packages",
            "platform": "ubuntu"
        }
    },
    "file_paths": {
        "configuration": [
            "/etc/rc%/%%",
            "/etc/init.d/%%",
            "/etc/systemd/%%",
            "/lib/systemd/%%",
            "/usr/lib/systemd/%%",
            "/etc/crontab",
            "/etc/cron.%/%%",
            "/var/spool/cron/%%",
            "/etc/nginx/%%",
            "/etc/ssh/sshd_config",
            "/etc/ld.so.preload",
            "/etc/ld.so.conf",
            "/etc/ld.so.conf.d/%%",
            "/etc/pam.d/%%",
            "/etc/update-motd.d/%%",
            "/etc/profile",
            "/etc/profile.d/%%",
            "/etc/shells",
            "/etc/bash.bashrc"
        ],
        "binaries": [
            "/usr/bin/%%",
            "/usr/sbin/%%",
            "/bin/%%",
            "/sbin/%%",
            "/usr/local/bin/%%",
            "/usr/local/sbin/%%"
        ],
        "users": [
            "/home/%%/%%",
            "/root/%%"
        ],
        "www": [
            "/var/www/%%"
        ],
        "tmp": [
            "/tmp/%%"
        ]
    },
    "packs": {
        "ossec-rootkit": "/opt/osquery/share/osquery/packs/ossec-rootkit.conf"
    }
}